# scripts/rollback/rollback-health-monitor.service
[Unit]
Description=TraceRTM Health Monitor and Automated Rollback Service
Documentation=https://github.com/tracertm/trace/docs/rollback
After=network.target prometheus.service
Wants=prometheus.service

[Service]
Type=simple
User=tracertm
Group=tracertm
WorkingDirectory=/opt/tracertm
Environment="CHECK_INTERVAL=30"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/opt/tracertm/scripts/rollback/health-monitor.sh start
ExecStop=/opt/tracertm/scripts/rollback/health-monitor.sh stop
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tracertm-health-monitor

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tracertm/.rollback-logs /opt/tracertm/.deployments

# Resource limits
MemoryMax=256M
CPUQuota=10%

[Install]
WantedBy=multi-user.target
